# chapter12/deploySetup/bookServer/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
nodaemon=true
user=root

[program:myserver]
command=./go/bin/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
