# chapter1/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:myserver]
command=/home/user/golab/bin/mirrorFinder
autostart=true
autorestart=true
redirect_stderr=true
